# files/etc/s6/nginx/finish
#!/bin/sh

kill "$(ps aux | grep '[n]ginx' | grep -v 's6' | awk '{print $1}')"
