# files/etc/s6/redis/finish
#!/bin/sh

kill "$(ps aux | grep '[r]edis' | grep -v 's6' | awk '{print $1}')"
