# files/etc/s6/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm -F --allow-to-run-as-root
