# files/etc/s6/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm -F --allow-to-run-as-root -y /etc/php81/deva/php-fpm.conf -c /etc/php81/deva/php.ini
