# files/etc/s6/mysql/finish
#!/bin/sh

kill "$(ps aux | grep '[m]ysql' | grep -v 's6' | awk '{print $1}')"
