# files/etc/s6/mysql/run
#!/bin/sh

exec /usr/bin/mysqld_safe --user=root --datadir='/var/lib/mysql'
