# rootfs/etc/supervisor.d/deva.ini
[supervisord]
nodaemon = true
user = root
pidfile = /run/supervisord.pid

[program:php-fpm]
command = /usr/sbin/php-fpm -F
numprocs = 1
autostart = true
autorestart = true
startsecs = 0

[program:nginx]
command = /usr/sbin/nginx
numprocs = 1
autostart = true
autorestart = true
startsecs = 0

[program:mysql]
command = /usr/bin/mysqld_safe --datadir='/var/lib/mysql'
numprocs = 1
autostart = true
autorestart = true
startsecs = 0
