# files/etc/supervisor.d/deva.ini
[supervisord]
nodaemon = true
user = root
pidfile = /run/supervisord.pid
logfile = /dev/null
logfile_maxbytes = 0

[program:php-fpm]
command = /usr/sbin/php-fpm -F --allow-to-run-as-root
numprocs = 1
autostart = true
autorestart = false
startsecs = 0

[program:nginx]
command = /usr/sbin/nginx
numprocs = 1
autostart = true
autorestart = false
startsecs = 0

[program:mysql]
command = /usr/bin/mysqld_safe --user=root --datadir='/var/lib/mysql'
numprocs = 1
autostart = true
autorestart = false
startsecs = 0

[program:redis]
command = /usr/bin/redis-server
numprocs = 1
autostart = false
autorestart = false
startsecs = 0
